# autostart/.config/autostart/shepherd.desktop
[Desktop Entry]
Name=Shepherd
GenericName=Shepherd
Exec=shepherd
Terminal=false
Type=Application
X-GNOME-Autostart-enabled=true
