# autostart/.config/autostart/syncthing.desktop
[Desktop Entry]
Name=Syncthing
GenericName=Syncthing
Exec=syncthing
Terminal=false
Type=Application
X-GNOME-Autostart-enabled=true
